# packages/kmod-6.1-nvidia/nvidia-driver.service
[Unit]
Description=Ensure the NVIDIA driver is loaded
Requires=nvidia-driver.path
# Disable manual restarts to prevent loading kernel modules
# that weren't linked by the running system
RefuseManualStart=true
RefuseManualStop=true

[Service]
Type=oneshot
ExecStart=/usr/bin/modinfo nvidia
RemainAfterExit=true
RestartSec=1
Restart=on-failure
StandardError=journal+console

[Install]
RequiredBy=preconfigured.target
