# packages/os/load-open-gpu-kernel-modules.service.in
[Unit]
Description=Load additional kernel modules
RequiresMountsFor=PREFIX/lib/modules PREFIX/src/kernels
After=link-open-gpu-kernel-modules.service
Requires=link-open-gpu-kernel-modules.service load-open-gpu-kernel-modules.path
ConditionPathExists=/run/nvidia/nvidia-open-gpu
# Disable manual restarts to prevent loading kernel modules
# that weren't linked by the running system
RefuseManualStart=true
RefuseManualStop=true

[Service]
Type=oneshot
ExecStart=/usr/bin/driverdog --modules-set nvidia-open-gpu load-modules
RemainAfterExit=true
StandardError=journal+console

